# web/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:php]
command=php-fpm -d variables_order="EGPCS" --nodaemonize
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:consul]
command=consul agent -retry-join "consul" -data-dir "/data/consul" -datacenter "sensu-clientless-test" -config-dir "/etc/consul.d/client"
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0